# crontab
*/5 * * * * cd /app && python backend/app/bot_runner.py >> /app/logs/cron.log 2>&1 